# systemd/user/onedrive-sync.service
[Unit]
Description=OneDrive Sync

[Service]
Type=oneshot
ExecStart=/usr/bin/onedrive --sync --verbose --single-directory="SharedDrive"

